# system/rc.local

mkdir -p /tmp/pike/inbox || true

exit 0
